# config/systemd/export-journal.timer
[Unit]
Description=Daily export of trading journal at 10:00 UTC

[Timer]
OnCalendar=*-*-* 10:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
